# deploy/systemd/i2pptt-backend.service.example
[Unit]
Description=i2pptt Backend Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/i2pptt
Environment="PATH=/path/to/i2pptt/venv/bin"
Environment="I2PPTT_ROOT_PATH=/i2pptt"
ExecStart=/path/to/i2pptt/venv/bin/python -m uvicorn web.backend.main:app --host 0.0.0.0 --port 8001 --workers 4
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

